# build/asus-smart-charge_0.1.0/stage/usr/lib/systemd/system/asus-smart-charge.service
[Unit]
Description=Apply Asus Smart Charge battery threshold
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/asus-smart-charge/asus-smart-charge-helper enforce

[Install]
WantedBy=multi-user.target
