# build/asus-smart-charge_0.1.0/AppDir/usr/lib/systemd/system/asus-smart-charge.timer
[Unit]
Description=Keep Asus Smart Charge threshold in sync

[Timer]
OnBootSec=20s
OnUnitActiveSec=1min
Persistent=true
Unit=asus-smart-charge.service

[Install]
WantedBy=timers.target
